# server configuration files/mysite.conf
[program:mysite]
command = /opt/mysite_env/bin/gunicorn_start
autostart=true
autorestart=true
stopasgroup=true
stderr_logfile=/var/log/long.err.log
stdout_logfile=/var/log/long.out.log
